# opensurveillancedb-ui/opensurveillancedb-caddy.service
[Unit]
Description=OpenSurveillanceDB Caddy TLS Reverse Proxy (Docker)
Requires=docker.service opensurveillancedb-ui.service
After=docker.service network-online.target opensurveillancedb-ui.service
Wants=network-online.target

[Service]
User=osdb
Group=osdb
Restart=always
RestartSec=5
TimeoutStartSec=0

# Set these for your deployment.
Environment=DOMAIN=opensurveillancedb.com
Environment=APP_UPSTREAM=127.0.0.1:8080

# Note: %h resolves to the home of the systemd manager (root), not User=osdb, so we use the absolute path instead.
ExecStart=/usr/bin/docker run --rm --name opensurveillancedb-caddy --network host --restart=no -v /home/osdb/opensurveillancedb/opensurveillancedb-ui/Caddyfile:/etc/caddy/Caddyfile:ro -v /home/osdb/.local/share/caddy/data:/data -v /home/osdb/.local/share/caddy/config:/config -e DOMAIN=${DOMAIN} -e APP_UPSTREAM=${APP_UPSTREAM} caddy:2-alpine
ExecStop=-/usr/bin/docker stop --time 20 opensurveillancedb-caddy

[Install]
WantedBy=multi-user.target
